# etc/service/mylar/run
#! /bin/sh
cd $APPDIR

# start mylar
exec /usr/bin/python Mylar.py -p 8090 --datadir=/config --nolaunch
